# ha-web-browser/rootfs/etc/services.d/addon/run
#!/usr/bin/with-contenv bashio

# Get configuration
RESOLUTION=$(bashio::config 'resolution')
HOMEPAGE=$(bashio::config 'homepage')

bashio::log.info "Starting Web Browser addon..."
bashio::log.info "Resolution: ${RESOLUTION}"
bashio::log.info "Homepage: ${HOMEPAGE}"

# Export environment variables for supervisor
export DISPLAY=:0
export RESOLUTION
export HOMEPAGE

# Start supervisor
exec /usr/bin/supervisord -c /etc/supervisor/conf.d/supervisord.conf
